# deploy/realhack-pilot.service
[Unit]
Description=RealHack Pilot FastAPI backend (uvicorn)
Documentation=https://github.com/ShaikshavaliKalluri/RealHackPilot
# Postgres lives on a remote VM (rcapaydbpgr001), so we only need the network
# to be up — no local postgresql.service to depend on.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=realhack
Group=realhack
WorkingDirectory=/opt/realhack-pilot/app/backend
EnvironmentFile=/opt/realhack-pilot/.env.production
ExecStart=/opt/realhack-pilot/app/backend/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2 --log-level info
Restart=on-failure
RestartSec=5

# Resource hints — conservative because this box also runs Pulse (.NET + Next.js)
# and Postgres. 16GB total; budgeting ~2GB hard ceiling for us is plenty for
# 48 teams' worth of data and a couple of concurrent LLM calls.
LimitNOFILE=4096
MemoryHigh=1G
MemoryMax=2G

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true

# Stream stdout/stderr to journalctl
StandardOutput=journal
StandardError=journal
SyslogIdentifier=realhack-pilot

[Install]
WantedBy=multi-user.target
